# dist/linux/debug-systemd/systemd.service
[Unit]
Description=Ghostty

[Service]
Type=dbus
BusName=com.mitchellh.ghostty-debug
ExecStart=ghostty
